# .config/systemd/user/restic-backup.service
[Unit]
Description=Restic backup

[Service]
ExecStart=systemd-inhibit $HOME/.local/bin/restic-backup